# config/exercise_supervisor.conf
[program:exercise]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn_config.py
directory=/var/www/exercise
autostart=true
autorestart=true